# etc/sytemd/system/lightdm.service.d/run-my-script-first.conf
[Service]
ExecStartPre=/home/matthieu/IR_Custom.sh
